# roles/win10/templates/win10-usb.service.j2
[Unit]
Description=Attach usb devices to vm

[Service]
ExecStart=/usr/bin/vfio-usb win10 --evdev-kbd "{{ kbd_device_path }}" {{ ' '.join(usb_device_ids) }}
Restart=always


