# roles/win10/templates/win10-usb.service.j2
[Unit]
Description=Attach usb devices to vm

[Service]
ExecStart=/usr/bin/vfio-usb win10 --evdev-kbd "/dev/input/by-id/uinput-persist-keyboard0" {{ ' '.join(usb_device_ids) }}
Restart=always


